# templates/consul-acl-bootstrap.service
[Unit]
Description=Register Nomad agent token in Consul ACL
After=consul.service
Requires=consul.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=ENROLL_ADDR=127.0.0.1:8443
Environment=ENROLL_CACERT=/etc/pigeon/identity/ca.pem
Environment=ENROLL_IDENTITY_DIR=/etc/pigeon/identity
ExecStart=/usr/local/bin/consul-acl-bootstrap

[Install]
WantedBy=multi-user.target
