# templates/luks-recovery-worker.service
[Unit]
Description=Add LUKS recovery passphrase to the boot volume
After=pigeon-template-reconcile.service
ConditionPathExists=/etc/pigeon/identity/cert.pem

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/pigeon/enroll.env
Environment=ENROLL_CACERT=/etc/pigeon/identity/ca.pem
Environment=ENROLL_IDENTITY_DIR=/etc/pigeon/identity
ExecStart=/usr/local/bin/luks-recovery

[Install]
WantedBy=multi-user.target
