# templates/nomad-cert.path
[Unit]
Description=Watch for Nomad TLS certificate from Vault PKI

[Path]
PathExists=/etc/nomad.d/certs/cert.pem
Unit=nomad.service

[Install]
WantedBy=multi-user.target
