# templates/pigeon-template-bootstrap.path
[Unit]
Description=Watch for pigeon enrollment secrets

[Path]
PathExists=/encrypted/pigeon/enroll.json
Unit=pigeon-template-bootstrap.service

[Install]
WantedBy=multi-user.target
