# templates/pigeon-template-secrets.path
[Unit]
Description=Watch for pigeon secrets file

[Path]
PathExists=/encrypted/pigeon/secrets.json
Unit=pigeon-template.service

[Install]
WantedBy=multi-user.target
