# templates/pigeon-template-bootstrap.service
[Unit]
Description=pigeon-template bootstrap config renderer (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/pigeon/bin/pigeon-template --once --config=/etc/pigeon/bootstrap.tmpl.hcl
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/etc/resolv.conf /etc/unbound/
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LimitCORE=0
