# templates/pigeon-fence.service
[Unit]
Description=pigeon-fence firewall manager
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/opt/pigeon/bin/pigeon-fence --config=/etc/pigeon/fence.hcl --config=/encrypted/pigeon/fence.d
Restart=on-failure
RestartSec=5

AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
ReadOnlyPaths=/etc/pigeon
ReadWritePaths=/encrypted/pigeon
MemoryMax=128M
LimitCORE=0

[Install]
WantedBy=multi-user.target
