# templates/pigeon-enroll-actions.service
[Unit]
Description=pigeon-enroll post-claim actions (vault-init, luks-recovery, etc.)
After=vault.service
Requires=vault.service
ConditionPathExists=!/encrypted/vault/init.json
ConditionFileNotEmpty=/encrypted/pigeon/enroll.hcl

[Service]
Type=oneshot
ExecStart=/opt/pigeon/bin/pigeon-enroll run-actions -config=/encrypted/pigeon/enroll.hcl
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
