# templates/pigeon-template-reconcile-worker.service
[Unit]
Description=pigeon-template reconciler
After=pigeon-template-bootstrap.service unbound.service pigeon-mesh.service
Wants=unbound.service pigeon-mesh.service
ConditionPathExists=/etc/pigeon/identity/cert.pem

[Service]
EnvironmentFile=/etc/pigeon/enroll.env
Environment=ENROLL_CACERT=/etc/pigeon/identity/ca.pem
Environment=ENROLL_IDENTITY_DIR=/etc/pigeon/identity
ExecStart=/opt/pigeon/bin/pigeon-template --config=/etc/pigeon/reconcile.hcl
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/etc/unbound/ /var/lib/pigeon/
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LimitCORE=0

[Install]
WantedBy=multi-user.target
