# templates/pigeon-enroll-actions.service
[Unit]
Description=pigeon-enroll post-claim actions (vault-init, luks-recovery, etc.)
After=vault.service
Requires=vault.service
ConditionPathExists=!/var/lib/pigeon/vault-init.json
ConditionFileNotEmpty=/etc/pigeon/enroll-server.hcl

[Service]
Type=oneshot
ExecStart=/opt/pigeon/bin/pigeon-enroll run-actions -config=/etc/pigeon/enroll-server.hcl
RemainAfterExit=yes
ProtectHome=yes
LimitCORE=0

[Install]
WantedBy=multi-user.target
