# templates/pigeon-template-bootstrap.path
[Unit]
Description=Trigger pigeon-template bootstrap when the identity cert is provisioned

[Path]
PathExists=/etc/pigeon/identity/cert.pem
Unit=pigeon-template-bootstrap.service

[Install]
WantedBy=multi-user.target
