# templates/pigeon-mesh.service
[Unit]
Description=pigeon-mesh WireGuard overlay
After=network-online.target
Wants=network-online.target
ConditionPathExists=/encrypted/pigeon/mesh.json

[Service]
ExecStart=/opt/pigeon/bin/pigeon-mesh --config=/encrypted/pigeon/mesh.json
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/pigeon-mesh /run/pigeon-mesh
ReadOnlyPaths=/encrypted
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
LimitCORE=0

[Install]
WantedBy=multi-user.target
