# templates/pigeon-template-reconcile-server.service
[Unit]
Description=pigeon-template reconciler
After=pigeon-enroll.service network-online.target unbound.service pigeon-mesh.service
Wants=network-online.target unbound.service pigeon-mesh.service
ConditionPathExists=/etc/pigeon/identity/cert.pem

[Service]
Type=notify
Environment=ENROLL_ADDR=127.0.0.1:8443
Environment=ENROLL_CACERT=/etc/pigeon/identity/ca.pem
Environment=ENROLL_IDENTITY_DIR=/etc/pigeon/identity
# Seed enroll.json before templates that reference ${file.enroll.*} try to
# render. After startup the exec source refresh_enroll_json keeps it fresh.
ExecStartPre=/opt/pigeon/bin/pigeon-enroll read template/enroll-server -o /var/lib/pigeon/enroll.json
ExecStart=/opt/pigeon/bin/pigeon-template --config=/etc/pigeon/reconcile.hcl
Restart=on-failure
RestartSec=5
WatchdogSec=30s
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/etc/pigeon/ /etc/vault.d/ /etc/consul.d/ /etc/nomad.d/ /etc/resolv.conf /etc/unbound/ /var/lib/pigeon/
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LimitCORE=0

[Install]
WantedBy=multi-user.target
