# templates/vault-agent-server.service
[Unit]
Description=Vault Agent (server TLS cert renewal from Vault PKI)
Wants=network-online.target
After=network-online.target pigeon-mesh.service consul.service vault.service
Requires=pigeon-mesh.service consul.service vault.service
ConditionPathExists=/etc/pigeon/certs/auth/cert.pem
ConditionPathExists=/etc/pigeon/certs/auth/key.pem
ConditionPathExists=/etc/pigeon/vault-agent.env

[Service]
EnvironmentFile=/etc/pigeon/vault-agent.env
ExecStart=/usr/bin/vault agent -config=/etc/pigeon/vault-agent.hcl
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitCORE=0
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/etc/consul.d/certs/ /etc/nomad.d/certs/ /etc/vault.d/certs/
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes

[Install]
WantedBy=multi-user.target
