# templates/vault-agent-server.service
[Unit]
Description=Vault Agent (server TLS cert renewal from Vault PKI)
Wants=network-online.target
After=network-online.target pigeon-mesh.service consul.service vault.service
Requires=pigeon-mesh.service consul.service
ConditionPathExists=/encrypted/tls/auth/cert.pem

[Service]
EnvironmentFile=/encrypted/pigeon/vault-agent.env
ExecStart=/usr/bin/vault agent -config=/etc/pigeon/vault-agent.hcl
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
