# templates/vault-agent-server.service
[Unit]
Description=Vault Agent (server TLS cert renewal from Vault PKI)
Wants=network-online.target
After=network-online.target pigeon-mesh.service consul.service vault.service
Requires=pigeon-mesh.service consul.service vault.service
ConditionPathExists=/encrypted/tls/auth/cert.pem
ConditionPathExists=/encrypted/tls/auth/key.pem
ConditionPathExists=/encrypted/pigeon/vault-agent.env

[Service]
EnvironmentFile=/encrypted/pigeon/vault-agent.env
ExecStart=/usr/bin/vault agent -config=/etc/pigeon/vault-agent.hcl
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitCORE=0
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes

[Install]
WantedBy=multi-user.target
